# crontab
* * * * * /usr/bin/python3 /app/Main.py
